# .config/autostart/random-wallpaper.desktop
[Desktop Entry]
Type=Application
Name=Random Wallpaper
Exec=/home/thomas/backgrounds/change-bg.sh
X-GNOME-Autostart-enabled=true
NoDisplay=true
